# DVDの中身/Ubuntu_etc_systemd_system/bluetooth_connect.service
[Unit]
Description=bluetooth service

[Service]
ExecStart=/bin/bash /home/tetsuhiro/28800.sh
type=simple

[Install]
WantedBy=multi-user.target

